# Dashboard/api_service.service
[Unit]
Description=Dashboard Api Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3:/usr/bin/bash /home/pi/Documents/RFID/Dashboard/run.sh
Restart=always

[Install]
WantedBy=multi-user.target